# docker/runit/nginx/run
#!/bin/sh 

set -e

nginx -g 'pid /var/run/nginx.pid; daemon off;'
